# files/cron.d/enforce_vpn_usage
# m     h       dom     mon     dow     user    command
@reboot                                 root    iptables -A OUTPUT ! -o tun0 -m owner --uid-owner debian-transmission -j REJECT --reject-with icmp-port-unreachable
